# configs/etc/s6/mosquitto/run
#!/bin/bash

# setup mosquitto
/opt/setup.sh

# start sshd
/usr/sbin/mosquitto -c /data/mqtt/config/mosquitto.conf
